# tool/xiuno_sphinx/jobs.cron
*/15 * * * * /usr/local/coreseek/bin/indexer -c /usr/local/coreseek/etc/csft.conf delta_thread --rotate > /dev/null 2>&1
5 4 * * * /usr/local/coreseek/bin/indexer -c /usr/local/coreseek/etc/csft.conf --merge bbs_thread delta_thread --rotate > /dev/null 2>&1