# tool/xiuno_sphinx/jobs.cron
*/15 * * * * /usr/local/coreseek/bin/indexer -c /usr/local/coreseek/etc/csft.conf delta_thread --rotate > /dev/null 2>&1
5 4 * * * /usr/local/coreseek/bin/indexer -c /usr/local/coreseek/etc/csft.conf --merge bbs_thread delta_thread --rotate > /dev/null 2>&1
5 3 * * * mysql -hlocalhost -uxiuno -pxiuno -e"use xiuno;REPLACE INTO bbs_framework_maxid SELECT 'sphinx_thread', MAX(tid) FROM bbs_thread";